# conf.d/gearmand.conf
[program:gearmand]
command=gearmand
priority=100
